# services/master_monitor.timer
[Unit]
Description=Run monitor every minute

[Timer]
OnUnitActiveSec=5min
OnBootSec=1min

[Install]
WantedBy=timers.target